# scripts/startup/j1msky-work-feed.service
[Unit]
Description=J1MSKY Work Feed
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/m1ndb0t/Desktop/J1MSKY
ExecStart=/usr/bin/python3 /home/m1ndb0t/Desktop/J1MSKY/dashboards/work-feed.py
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
